# setup/docker_config/upstart/redis.conf
description "redis"

start on startup
stop on shutdown

setuid redis
setgid redis

respawn

exec start-stop-daemon --start --exec /usr/local/bin/redis-server /etc/redis/redis.conf

expect fork
